# src/archlinux_management/resources/99-mute-yeti-mic-loopback.rules
ACTION=="add", SUBSYSTEM=="sound", ATTRS{idVendor}=="b58e", ATTRS{idProduct}=="9e84", RUN+="/usr/local/bin/mute-yeti-mic-loopback.sh"
